# supervisor-wvg-lock.conf
[program:wvg-lock]
directory=/opt/wvg/pi-nfc-lock/
command=python pi-nfc-lock.py
autorestart=true
autostart=true
